# package/root/etc/update-motd.d/05-figlet
#!/bin/sh

if [ $(which figlet) ]; then
figlet $(hostname)
fi
